# supervisor-conf.d/sonos-yamaha-monitor_family-room.conf
[program:sonos-yamaha-monitor_family-room]
command=/usr/local/bin/pew in sonos-yamaha-monitor ./sonos-yamaha-monitor.py --zone family-room
numprocs=1
directory=/home/wylie/src/axontrust/sonos-yamaha-monitor
autostart=true                ; start at supervisord start (default: true)
autorestart=unexpected        ; whether/when to restart (default: unexpected)
startsecs=1                   ; number of secs prog must stay running (def. 1)
startretries=3                ; max # of serial start failures (default 3)
exitcodes=0,2                 ; 'expected' exit codes for process (default 0,2)
stopsignal=QUIT               ; signal used to kill process (default TERM)
stopwaitsecs=10               ; max num secs to wait b4 SIGKILL (default 10)
user=wylie                    ; setuid to this UNIX account to run the program
redirect_stderr=true          ; redirect proc stderr to stdout (default false)
stdout_logfile=/var/log/sonos-yamaha-monitor/family-room.log
stdout_logfile_maxbytes=10MB   ; max # logfile bytes b4 rotation (default 50MB)
stdout_logfile_backups=10     ; # of stdout logfile backups (default 10)
environment=HOME=/home/wylie,PYTHONUNBUFFERED=True
stopasgroup=true        ; send stop signal to the UNIX process group (default false)
killasgroup=true        ; SIGKILL the UNIX process group (def false)
